# cygnum.service
# sudo ln -s /home/ubuntu/cygnum/cygnum.service /etc/systemd/system/cygnum.service
[Unit]
Description=Gunicorn instance to serve cygnum
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/home/ubuntu/cygnum
Environment="PATH=/home/ubuntu/cygnum/env/bin"
Environment="CYGNUM_CONFIG=/home/ubuntu/cygnum/instance/config.py"
ExecStart=/home/ubuntu/cygnum/env/bin/gunicorn --config gunicorn.conf.py wsgi:app

[Install]
WantedBy=multi-user.target
